# selimcam.service
[Unit]
Description=SelimCam - Raspberry Pi Camera Application
After=network.target multi-user.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/FINALMAINCAMMM

# No terminal output - completely silent
StandardOutput=null
StandardError=null

# Quiet mode - no debug spam
Environment="SELIMCAM_QUIET=true"
Environment="DISPLAY=:0"

# Run the app
ExecStartPre=/home/pi/FINALMAINCAMMM/show_boot_screen.sh
ExecStart=/home/pi/FINALMAINCAMMM/.venv/bin/python main.py

# Restart on crash with delay
Restart=always
RestartSec=10

# Kill process after 10 seconds on stop
KillMode=process
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target

